# contrib/init.d/rc.gentoo.openpbx
#!/sbin/runscript

depend() {
	need net logger
}

start() {
	ebegin "Starting CallWeaver"
	/usr/sbin/openpbx
	eend $? "Failed to start CallWeaver"
}

stop() {
	ebegin "Stopping CallWeaver"
	kill $(cat /var/run/openpbx.pid)
	eend $? "Failed to stop CallWeaver"
}
