# system/rc.callweaver.gentoo.in
#!/sbin/runscript

depend() {
	need net logger
}

start() {
	ebegin "Starting CallWeaver"
	'@cwexecdir@/callweaver'
	eend $? "Failed to start CallWeaver"
}

stop() {
	ebegin "Stopping CallWeaver"
	kill $(cat '@cwpidfile@')
	eend $? "Failed to stop CallWeaver"
}
