# system/rc.callweaver.redhat.in
#!/bin/sh
#
# callweaver    Starts, Stops and Reloads CallWeaver.
#
# chkconfig: 2345 40 60
# description: CallWeaver PBX server
# processname: callweaver
# pidfile: @cwpidfile@

# Source function library.
. /etc/rc.d/init.d/functions

[ -x '@cwexecdir@/callweaver' ] || exit 0
[ -d /etc/callweaver ] || exit 0

# Full path to callweaver binary
DAEMON='@cwexecdir@/callweaver'

# Full path to safe_callweaver script
SAFE_CALLWEAVER='@cwutilsdir@/safe_callweaver'

RETVAL=0

start() {
	# Start daemons.
	echo -n 'Starting callweaver: '
	daemon "${SAFE_CALLWEAVER:-$DAEMON}"
	RETVAL=$?
	[ $RETVAL -eq 0 ] && touch /var/lock/subsys/callweaver
	echo
	return $RETVAL
}

stop() {
	# Stop daemons.
	echo -n 'Shutting down callweaver: '
	killproc callweaver
	RETVAL=$?
	[ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/callweaver
	echo
	return $RETVAL
}

restart() {
	stop
	start
}

reload() {
	killproc callweaver -HUP
}

# See how we were called.
case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  restart)
	restart
	;;
  reload)
	reload
	;;
  condrestart)
	[ -f /var/lock/subsys/callweaver ] && restart || :
	;;
  status)
	status callweaver
	;;
  *)
	echo "Usage: ${0##*/} {start|stop|restart|reload|condrestart|status}"
	exit 1
esac

exit $?
